# deploy/docker.service.d/wait-for-tailscale.conf
# systemd drop-in: make Docker daemon wait for Tailscale at boot.
#
# Without this, on jib-jab (Debian 13 with auto-WiFi via NetworkManager and
# tailscaled running as a separate systemd service), Docker can start
# BEFORE tailscaled has finished connecting and assigning the tailnet IP
# to the host. When dockerd then tries to bring up containers via the
# restart policy, the Postgres container's port binding
# (`${TAILSCALE_IP}:5432:5432` in `docker-compose.primary.yml`) fails
# because the tailnet IP isn't on any interface yet. containerd kills the
# container, and the stack never recovers without manual intervention.
#
# This drop-in ensures Docker daemon (and therefore everything it manages
# at runtime, including `restart: unless-stopped` containers and the
# hrserv.service systemd unit's `dc up -d`) only starts AFTER:
#   - network-online.target — host network is up
#   - tailscaled.service     — Tailscale daemon is up AND tailnet IP is assigned
#
# Install on a new node (one-time, requires root):
#   sudo mkdir -p /etc/systemd/system/docker.service.d
#   sudo cp /opt/hrserv/deploy/docker.service.d/wait-for-tailscale.conf \
#       /etc/systemd/system/docker.service.d/
#   sudo systemctl daemon-reload
# The change takes effect on the next reboot.
#
# `Wants=` is intentionally soft — if tailscaled is somehow not installed,
# this becomes a no-op rather than blocking Docker startup forever.

[Unit]
After=tailscaled.service network-online.target
Wants=tailscaled.service network-online.target
