# deploy/docker.service.d/wait-for-tailscale.conf
# systemd drop-in: make Docker daemon wait for Tailscale at boot.
#
# Without this, on jib-jab (Debian 13 with auto-WiFi via NetworkManager and
# tailscaled running as a separate systemd service), Docker can start
# BEFORE tailscaled has finished connecting and assigning the tailnet IP
# to the host. When dockerd then tries to bring up containers via the
# restart policy, the Postgres container's port binding
# (`${TAILSCALE_IP}:5432:5432` in `docker-compose.primary.yml`) fails with
# "cannot assign requested address" because the tailnet IP isn't on any
# interface yet. containerd kills the container, hrserv.service exits with
# status 1, and the stack stays down until manual `systemctl restart hrserv`.
# Observed on jib-jab 2026-05-16: tailscaled needed ~10-15s after boot to
# assign 100.91.182.4 to tailscale0, but hrserv.service fired at ~10s
# (After=tailscaled.service was satisfied as soon as the daemon started,
# NOT when it finished negotiating the IP). 16-minute outage until manual
# restart succeeded.
#
# Three-layer fix:
#
# 1. `After=tailscaled.service` — orders Docker to start after tailscaled
#    has *started*. Necessary but not sufficient: tailscaled being "active"
#    just means the daemon is running, not that tailscale0 has its IP.
#
# 2. `ExecStartPre=/usr/bin/tailscale wait` — blocks until tailscale0 is
#    actually bindable. This is Tailscale's own command for exactly this
#    use case (added in Tailscale 1.32+). It exits 0 once the interface
#    has its IP, or stays blocked until it does.
#
# 3. `TimeoutStartSec=120` — IMPORTANT. Upstream docker.service ships with
#    TimeoutStartSec=0, which systemd interprets as INFINITY. Without this
#    override, if tailscaled never assigns an IP (auth key expired, control
#    plane unreachable, tailnet ACL change, mass node revocation), the
#    `tailscale wait` ExecStartPre blocks forever and Docker silently stays
#    down. Setting 120s bounds the wait: 10-15s is the observed normal
#    case on jib-jab, so 120s is a generous-but-bounded headroom. If we
#    hit it, the unit fails with a clear journal entry — paging/monitoring
#    can catch that, whereas an infinite hang is invisible.
#
# This drop-in ensures Docker daemon (and therefore everything it manages
# at runtime, including `restart: unless-stopped` containers and the
# hrserv.service systemd unit's `dc up -d`) only starts AFTER:
#   - network-online.target — host network is up
#   - tailscaled.service     — Tailscale daemon has started
#   - `tailscale wait`        — tailscale0 has its IP and is bindable
#
# SCOPE: install this drop-in ONLY on hosts where docker-compose binds a
# tailnet IP (currently jib-jab via `docker-compose.primary.yml`, and the
# eventual replica via `docker-compose.replica.yml`). Local dev machines
# running `docker-compose.test.yml` (which binds 127.0.0.1) should NOT
# install this — it would unnecessarily couple dev Docker startup to
# Tailscale availability.
#
# Install on a new node (one-time, requires root). PRECONDITION: verify
# `which tailscale` returns `/usr/bin/tailscale` (apt-installed Tailscale).
# If installed via static binary or homebrew, the path differs and this
# drop-in must be updated accordingly — systemd requires absolute paths.
#
#   which tailscale                                  # must print /usr/bin/tailscale
#   tailscale version                                # must be >= 1.32 (for `wait`)
#   sudo mkdir -p /etc/systemd/system/docker.service.d
#   sudo cp /opt/hrserv/deploy/docker.service.d/wait-for-tailscale.conf \
#       /etc/systemd/system/docker.service.d/
#   sudo systemctl daemon-reload
#   systemd-analyze verify /etc/systemd/system/docker.service              # sanity check
# The change takes effect on the next reboot.
#
# `Wants=` is intentionally soft on the systemd dependency side — if
# tailscaled is somehow not installed, the ordering becomes a no-op. But
# `tailscale wait` is unconditional in ExecStartPre, so if tailscale is
# uninstalled while this drop-in is active, Docker startup will fail with
# `status=203/EXEC` (clear error, easy to diagnose). That's the right
# behavior: don't silently mask a misconfig.

[Unit]
After=tailscaled.service network-online.target
Wants=tailscaled.service network-online.target

[Service]
TimeoutStartSec=120
ExecStartPre=/usr/bin/tailscale wait
